# fs/etc/services.d/libvirtd/run
#!/bin/bash
mkdir -p /var/run/libvirt
mkdir -p /var/log/libvirt
# Lancer libvirtd en avant-plan pour s6
exec /usr/sbin/libvirtd --listen
